# qr_bot.service
[Service]
WorkingDirectory= /home/tg_bot/
User=root
ExecStart=/usr/bin/python3 tg_bot.py

[Install]
WantedBy=multi-user.target
EOF